# snap-scheduler.service
[Unit]
Description=Snapshot Scheduler Service
After=network.target

[Service]
Type=simple
ExecStart=/opt/stack/snapshot-scheduler/venv/bin/python /opt/stack/snapshot-scheduler/app/main.py
WorkingDirectory=/opt/stack/snapshot-scheduler/app
Restart=always
User=stack
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
