# systemd/cloudwatcher.service
[Unit]
Description=CloudWatcher Dashboard and Tagger
After=network.target

[Service]
Type=simple
Environment="OS_CLOUDS_YAML=/opt/stack/cloudwatcher/config/clouds.yaml"
Environment="OS_AUTH_URL=http://127.0.0.1:5000"
Environment="OS_USERNAME=demo"
Environment="OS_PASSWORD=secret"
Environment="OS_PROJECT_NAME=default"
ExecStart=/opt/stack/cloudwatcher/venv/bin/python3 /opt/stack/cloudwatcher/bin/quota_dashboard.py
ExecStartPost=/opt/stack/cloudwatcher/venv/bin/python3 /opt/stack/cloudwatcher/bin/weather_tagger.py
RemainAfterExit=true
StandardOutput=journal
StandardError=journal
User=stack

[Install]
WantedBy=multi-user.target
