# systemd/cloudcost_idle.service
[Unit]
Description=CloudCost-Idle Service
After=network.target

[Service]
User=stack
WorkingDirectory=/opt/stack/cloudcost-idle
Environment="PYTHONPATH=/opt/stack/cloudcost-idle"
ExecStart=/opt/stack/cloudcost-idle/venv/bin/python3 -m app.dashboard
Restart=always

[Install]
WantedBy=multi-user.target
