# systemd/cloudwatcher_dashboard.service
[Unit]
Description=CloudWatcher Dashboard
After=network.target

[Service]
Type=simple
Environment=OS_AUTH_URL=http://10.0.2.15/identity
Environment=OS_AUTH_TYPE=password
Environment=OS_USERNAME=admin
Environment=OS_PASSWORD=secret
Environment=OS_PROJECT_NAME=admin
Environment=OS_USER_DOMAIN_ID=default
Environment=OS_PROJECT_DOMAIN_ID=default
Environment=OS_REGION_NAME=RegionOne
Environment=OS_VOLUME_API_VERSION=3
ExecStart=/opt/stack/cloudwatcher/venv/bin/python3 /opt/stack/cloudwatcher/bin/quota_dashboard.py
Restart=always
User=stack

[Install]
WantedBy=multi-user.target

