# systemd/cloudcost_idle.service
[Unit]
Description=CloudCost Idle Plugin
After=network.target

[Service]
Type=simple
ExecStart=/opt/stack/Cloudcost-Idle/venv/bin/python3 /opt/stack/Cloudcost-Idle/app/dashboard.py
WorkingDirectory=/opt/stack/Cloudcost-Idle
EnvironmentFile=/opt/stack/Cloudcost-Idle/config/settings
Restart=always
User=stack

[Install]
WantedBy=multi-user.target
