# systemd/cloudwatcher.service
[Unit]
Description=CloudWatcher Dashboard and Tagger
After=network.target

[Service]
Type=simple
ExecStart=/opt/stack/cloudwatcher/venv/bin/python3 /opt/stack/cloudwatcher/bin/quota_dashboard.py
ExecStartPost=/opt/stack/cloudwatcher/venv/bin/python3 /opt/stack/cloudwatcher/bin/weather_tagger.py
RemainAfterExit=true
StandardOutput=journal
StandardError=journal
User=stack

[Install]
WantedBy=multi-user.target
