# systemd/cloudwatcher_dashboard.service
[Unit]
Description=CloudWatcher Dashboard
After=network.target

[Service]
Type=simple
Environment="OS_CLOUDS_YAML=/opt/stack/cloudwatcher/config/clouds.yaml"
ExecStart=/opt/stack/cloudwatcher/venv/bin/python3 /opt/stack/cloudwatcher/bin/quota_dashboard.py
Restart=on-failure
User=stack

[Install]
WantedBy=multi-user.target
