# systemd/cloudwatcher_tagger.service
[Unit]
Description=CloudWatcher Weather Tagger
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash -c "source /opt/stack/devstack/openrc admin admin && /opt/stack/cloudwatcher/venv/bin/python3 /opt/stack/cloudwatcher/bin/weather_tagger.py"
Restart=always
User=stack

[Install]
WantedBy=multi-user.target
